# crontab
#MAILTO=""
*/1 6-20 * * * /home/pi/dev/piTimelapseLinus/shoot.sh # this one actually runs until 20:59
0      0 * * * /home/pi/dev/piTimelapseLinus/cleanup.sh //delete all images older that one week
